# etc/systemd/system/tweetcontrol.service
[Unit]
Description = Check twitter for commands

[Service]
ExecStart = /usr/local/bin/tweetcontrol.sh

[Install]
WantedBy = mult-iuser.target
